# ah-shit-switch.service
[Unit]
Description=ah-shit-switch — wrong-layout corrector for Hyprland
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ah-shit-switch
Restart=on-failure
RestartSec=3

[Install]
WantedBy=graphical-session.target
